# Iris_fastapi_supervisor.conf
[program:webservice]
command=/home/ubuntu/venvs/Iris-fastapi/bin/gunicorn -k uvicorn.workers.UvicornWorker webapp:app -b 0.0.0.0:8080
directory=/home/ubuntu/opt/Iris_fastapi
user=ubuntu
